# udev_rules/monitor-detect.rules
ACTION=="change", RUN+="/home/kovas/bin/setup-monitors.bash"
